# skill-sync-user.service
[Unit]
Description=Auto sync OpenCode skills to Google Drive (user service)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/skill-sync.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
